# scripts/plists/com.wellgreens.pos-sync-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wellgreens.pos-sync-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO_DIR__/scripts/sync-weekly.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Monday at 2:00 AM -->
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/logs/sync-weekly.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/logs/sync-weekly.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
